# tmp/droid-config.files

/lib/udev/rules.d/997-droid-bsg.rules
/lib/udev/rules.d/997-droid-drm.rules
/lib/udev/rules.d/998-droid-system.rules
/var/lib/environment/ofono/noplugin.conf
/var/lib/environment/nemo/60-multimedia.conf
/var/lib/environment/compositor/droid-hal-device.conf
/usr/lib/sysctl.d/ipv6-settings.conf
/usr/lib/sysctl.d/ipv4-settings.conf
/usr/lib/tmpfiles.d/droid-hal.conf
/usr/lib/systemd/user/jolla-startupwizard-pre-user-session.service.d/50-compositor.conf
/usr/lib/systemd/user/jolla-actdead-charging.service.d/50-compositor.conf
/usr/lib/systemd/user/lipstick.service.d/50-compositor.conf
/usr/lib/systemd/system/droid-hal-init.service
/usr/lib/systemd/system/droid-reclaim-memory.timer
/usr/lib/systemd/system/droid-late-start.service
/usr/lib/systemd/system/graphical.target.wants/droid-hal-init.service
/usr/lib/systemd/system/graphical.target.wants/droid-late-start.service
/usr/lib/systemd/system/actdead.target.wants/droid-hal-init.service
/usr/lib/systemd/system/droid-reclaim-memory.service
/usr/lib/systemd/system/timers.target.wants/droid-reclaim-memory.timer
/usr/bin/droid/droid-init-done.sh
/usr/bin/droid/droid-hal-startup.sh
/usr/bin/droid/droid-hal-shutdown.sh
/usr/bin/droid/get_partition.sh
/usr/bin/droid/droid-reclaim-memory.sh
/usr/bin/droid/droid-load-firmware.sh
/etc/gps.conf
/etc/izat.conf
/etc/media_codecs.xml
/etc/sec_config
/etc/sap.conf
/etc/media_codecs_performance.xml
/etc/ofono/binder.d/dual-sim.conf
/etc/ofono/binder.conf
/etc/acdbdata
/etc/xtra_root_cert.pem
/etc/xtwifi.conf
/etc/nfcd/defaults.d/off-by-default.conf
/etc/selinux/config
/etc/media_codecs_google_audio.xml
/etc/modprobe.d/blacklist-evbug.conf
/etc/media_profiles.xml
/etc/media_codecs_google_video.xml
/etc/xdg/QtProject/QtQuick2.conf
/etc/xdg/QtProject/QPlatformTheme.conf
/etc/xdg/QtProject/Sensors.conf
/etc/obexd/excludes/bluetooth:pcsuite
/etc/obexd/pbap.conf
/etc/obexd/opp.conf
/etc/obexd/root
/etc/obexd/noplugins/mas
/etc/obexd/noplugins/irmc
/etc/obexd/noplugins/ftp
/etc/obexd/noplugins/syncevolution
/etc/lowi.conf
/etc/media_codecs_ffmpeg.xml
/etc/sensorfw/primaryuse.conf
/etc/mixer_paths.xml
/etc/media_codecs_google_telephony.xml
/etc/usb-moded/dyn-modes/mtp_mode-droid.ini
/etc/audio_platform_info.xml
/usr/lib/systemd/system/droid-hal-init.service
/usr/lib/systemd/system/wait_for_keymaster.service
/usr/lib/systemd/system/droid-bootctl.service
/usr/lib/systemd/system/graphical.target.wants/droid-bootctl.service
/usr/lib/systemd/system/multi-user.target.wants/wait_for_keymaster.service
/usr/share/ngfd/plugins.d/60-profile.ini
/usr/libexec/droid-hybris/system/etc/init/ecclist.rc
/usr/libexec/droid-hybris/system/etc/init/disabled_services.rc
/usr/bin/droid/wait_for_keymaster
/usr/bin/droid/droid-bootctl.sh
/boot/extracting-README.txt
/boot/flashing-README.txt
/etc/mce/60-hybris-features.ini
/etc/mce/60-fingerprint-sensor.conf
/etc/mce/60-proximity-sensor.conf
/etc/sailfish-fpd/50-settings.ini
/etc/gbinder.conf
/etc/gst-droid/gstdroidcodec.conf
/etc/ofono/ril_subscription.conf
/etc/ofono/ril_subscription.d/dual-sim.conf
/var/lib/environment/compositor/droid-hal-device.conf~
/var/lib/environment/compositor/droid-hal-device.conf
/usr/lib/systemd/system/system_ext.mount
/usr/lib/systemd/system/system_root.mount
/usr/lib/systemd/system/wifisetup.service
/usr/lib/systemd/system/product.mount
/usr/lib/systemd/system/dmsetup.service
/usr/lib/systemd/system/vendor-lib64-hw-audio.primary.default.so.mount
/usr/lib/systemd/system/network.target.wants/wifisetup.service
/usr/lib/systemd/system/odm.mount
/usr/lib/systemd/system/vendor_dlkm.mount
/usr/lib/systemd/system/vendor.mount
/usr/lib/systemd/system/local-fs.target.wants/system_ext.mount
/usr/lib/systemd/system/local-fs.target.wants/system_root.mount
/usr/lib/systemd/system/local-fs.target.wants/product.mount
/usr/lib/systemd/system/local-fs.target.wants/dmsetup.service
/usr/lib/systemd/system/local-fs.target.wants/vendor-lib64-hw-audio.primary.default.so.mount
/usr/lib/systemd/system/local-fs.target.wants/odm.mount
/usr/lib/systemd/system/local-fs.target.wants/vendor_dlkm.mount
/usr/lib/systemd/system/local-fs.target.wants/vendor.mount
/usr/lib/systemd/system/local-fs.target.wants/system.mount
/usr/lib/systemd/system/system.mount
/usr/share/csd/settings.d/10-hwsettings.ini
/usr/libexec/droid-hybris/system/etc/init/hwservicemanager.rc
/usr/libexec/droid-hybris/system/etc/init/ofono-ril-ecclists.rc
/usr/libexec/droid-hybris/system/etc/init/vndservicemanager.rc
/usr/libexec/droid-hybris/system/etc/init/disabled_services.rc
/usr/bin/droid/dmsetup.sh
/usr/bin/droid/droid-hal-early-init.sh
/linkerconfig/ld.config.txt
/system/lib/libandroidicu.so
/system/lib/libnativehelper.so
/system/lib/libicui18n.so
/system/lib/libicuuc.so
/system/lib64/libandroidicu.so
/system/lib64/libicui18n.so
/system/lib64/libicuuc.so
/etc/modules-load.d/pdx225.conf~
/etc/modules-load.d/pdx225.conf
/etc/binder.d/pdx225.conf
/etc/mce/70-battery-pdx225.ini
/etc/selinux/config
/etc/selinux/minimum/policy/policy.30
/etc/selinux/minimum/contexts/files/file_contexts
/etc/selinux/minimum/contexts/dbus_contexts
/etc/task_profiles
/etc/task_profiles.json
/usr/share/ssu/features.d/adaptation-community.ini
