# delete_file_xqcc72.list
/etc/sysconfig/dummy_netd
/usr/lib/systemd/system/wait_for_keymint.service